# yify-pop/src/supervisord.conf
[supervisord]
nodaemon=true
[program:yify-pop]
command=/opt/yify-pop/geddy
