# yify-pop/src/supervisord.conf
[supervisord]
nodaemon=true
[program:yify-pop]
command=/usr/local/bin/geddy --geddy-root /opt/yify-pop
