# yify-pop/src/supervisord.conf
[supervisord]
nodaemon=true
[program:yify-pop]
command=/usr/bin/geddy
