# sonarr/src/supervisord.conf
[supervisord]
nodaemon=true

[program:sonarr]
user=nobody
command=mono /opt/NzbDrone/NzbDrone.exe --no-browser -data=/config/sonarr

