# websync/services/websync/run
#!/bin/bash
cd /app/websync/dist
/sbin/setuser abc node server.js
