# minetest/src/supervisord.conf
[supervisord]
nodaemon=true
[program:minetest-server]
user=minetest
command=/usr/local/bin/minetestserver --gameid minimal --config /home/minetest/.minetest/minetest.conf
