# ghost/src/supervisord.conf
[supervisord]
nodaemon=true
[program:ghost]
user=nobody
directory=/opt/ghost
command=npm start
