# hublin/src/supervisord.conf
[supervisord]
nodaemon=true
[program:mongo]
user=root
command=/usr/bin/mongod
[program:hublin]
directory=/opt/meetings
user=root
command=npm start
