# pacs/src/supervisord.conf
[supervisord]
nodaemon=true
[program:mysql-server]
user=nobody
command=/usr/bin/mysqld_safe --skip-syslog
[program:pacs-server]
user=nobody
command="%(ENV_DCM4CHEE_HOME)s%(ENV_DCM_RUN)s"
