# pacs/src/supervisord.conf
[supervisord]
nodaemon=true
[program:mysql-server]
user=nobody
command=/usr/bin/mysqld_safe --skip-syslog
[program:pacs-server]
user=nobody
command=/config/dcm4chee/dcm4chee-mysql/bin/run.sh
