# daapd-server/src/supervisord.conf
[supervisord]
nodaemon=true
[program:avahi-daemon]
command=avahi-daemon -c /config/avahi-daemon.conf --no-chroot
[program:daapd-server]
command=/usr/sbin/forked-daapd -c /config/forked-daapd.conf
