# daapd-server/src/supervisord.conf
[supervisord]
nodaemon=true
[program:daapd-server]
command=/usr/sbin/forked-daapd -c /config/forked-daapd.conf
