# znc/src/supervisord.conf
[supervisord]
nodaemon=true
[program:znc]
user=nobody
command=/usr/local/bin/znc -d /config
