# pynab/src/supervisor-files/nginx-supervisord.conf
[supervisord]
nodaemon=true
[program:uwsgi-api]
command = /usr/local/bin/uwsgi --ini /root/pynab.ini

[program:nginx-api]
command = /usr/sbin/nginx
