# gerrit/src/supervisord.conf
[supervisord]
nodaemon=true
[program:gerrit]
user=nobody
command=/opt/gerrit/bin/gerrit.sh start
