# gerrit/src/supervisord.conf
[supervisord]
nodaemon=true

[program:gerrit]
command=/bin/bash -c "/home/gerrit/gerrit/bin/gerrit.sh start"