# bright-game/crontab
*/10 * * * * php -q /var/www/bgpanel/admin/cron.php > /dev/null 2>&1
