# amule/src/supervisord.conf
[supervisord]
nodaemon=true
[program:amule-server]
user=nobody
command=/opt/amule/amule/bin/bin/amuled -o -c /config -f --use-amuleweb=/opt/amule/amule/bin/bin/amuleweb
