# tftp-server/src/supervisord.conf
[supervisord]
nodaemon=true
[program:tftp-server]
user=root
command=/usr/sbin/in.tftpd --foreground --user tftp --address 0.0.0.0:69 --secure /images
