# minetest/src/service/minetest/run
#!/bin/bash
exec /sbin/setuser minetest /opt/minetest/bin/minetestserver --config /config/main-config/minetest.conf --port 30000
